# engine/scout/defaults/com.scout.heartbeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.scout.heartbeat.plist — runs ~/Scout/scripts/heartbeat.sh every 30 min.

  Installed by `scoutctl schedule install-heartbeat-plist`, which fills
  __USER_HOME__ at install time. Mirrors the schedule-tick plist pattern.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.scout.heartbeat</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__USER_HOME__/Scout/scripts/heartbeat.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__USER_HOME__/Scout</string>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__USER_HOME__/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>__USER_HOME__</string>
    </dict>

    <key>StandardOutPath</key>
    <string>__USER_HOME__/Scout/.scout-logs/launchd-heartbeat.log</string>
    <key>StandardErrorPath</key>
    <string>__USER_HOME__/Scout/.scout-logs/launchd-heartbeat.log</string>
</dict>
</plist>
